# roles/forgejo/templates/com.forgejo.backup.plist.j2
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.forgejo.backup</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>{{ ansible_env.HOME }}/.forgejo/backup.sh</string>
    </array>

    <key>StartInterval</key>
    <integer>43200</integer> <!-- Run every 12 hours -->

    <key>StandardOutPath</key>
    <string>{{ ansible_env.HOME }}/Library/Logs/forgejo-backup.log</string>

    <key>StandardErrorPath</key>
    <string>{{ ansible_env.HOME }}/Library/Logs/forgejo-backup.error.log</string>
</dict>
</plist>
